# deploy/continuous-listener.service
[Unit]
Description=Voice Assistant Continuous Listener
After=network.target

[Service]
Type=simple
User=jirka
WorkingDirectory=/home/jirka/virtual-assistant/continuous-listener
Environment="DOTNET_ROOT=/home/jirka/.dotnet"
Environment="XDG_RUNTIME_DIR=/run/user/1000"
Environment="PULSE_SERVER=unix:/run/user/1000/pulse/native"
ExecStart=/home/jirka/.dotnet/dotnet /home/jirka/virtual-assistant/continuous-listener/ContinuousListener.dll
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=continuous-listener

[Install]
WantedBy=multi-user.target
