# deploy/push-to-talk-dictation.service
[Unit]
Description=Push-to-Talk Dictation Service
After=graphical-session.target

[Service]
Type=simple
User=jirka
Environment=DISPLAY=:0
Environment=XDG_SESSION_TYPE=wayland
Environment=WAYLAND_DISPLAY=wayland-0
Environment=XAUTHORITY=/run/user/1000/.mutter-Xwaylandauth.0KYPG3
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=PULSE_SERVER=unix:/run/user/1000/pulse/native
EnvironmentFile=%h/.config/systemd/user/virtual-assistant.env
Environment=DOTNET_ROOT=%h/.dotnet
Environment=LD_LIBRARY_PATH=${VIRTUAL_ASSISTANT_BASE}/push-to-talk/runtimes/cuda/linux-x64:${VIRTUAL_ASSISTANT_BASE}/push-to-talk/runtimes/linux-x64:/usr/local/cuda-13.0/lib64
WorkingDirectory=${VIRTUAL_ASSISTANT_BASE}/push-to-talk
ExecStart=%h/.dotnet/dotnet ${VIRTUAL_ASSISTANT_BASE}/push-to-talk/VirtualAssistant.PushToTalk.Service.dll
Restart=on-failure
RestartSec=5
TimeoutStopSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
